# infra/systemd/tidepool.service
[Unit]
Description=Tidepool control plane (reconcile loop)
Documentation=https://github.com/0x63616c/tidepool
# Wait for the network and the state volume's mount unit (systemd derives
# `mnt-tidepool.mount` from the /mnt/tidepool mount point) before starting.
After=network-online.target mnt-tidepool.mount
Wants=network-online.target
# Stay inert until the operator delivers the age master key out-of-band. The box
# can boot, enable this unit, and still never run the control plane until the key
# is in place — secrets stay home (tenet 9).
ConditionPathExists=/root/.tidepool/bootstrap/age-mainbox.key

[Service]
Type=simple
WorkingDirectory=/opt/tidepool
# JIT-decrypt the sops bundle into the on-box runtime layout (incl. this env
# file) before the control plane starts. Optional load (`-`) so the very first
# activation does not hard-fail before ExecStartPre has generated the file; the
# always-restart loop re-reads it on the next cycle.
EnvironmentFile=-/etc/tidepool/env
ExecStartPre=/usr/local/bin/materialize-secrets.sh
ExecStart=/root/.bun/bin/bun /opt/tidepool/src/cli.ts run --watch
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
